# supervisor/registry.conf
[program:registry]
directory=/opt/dockyard-v2/scripts
command=python3 run_registry.py
