# supervisor/registry.conf
[program:registry]
directory=/opt/dockyard-v2/scripts
command=python run_registry.py
