# supervisor/registry.conf
[program:registry]
directory=/opt/dockyard-v2/scripts/run_registry
command=python run.py

[program:cron]
command=cron -f
